# systemd-templates/gen-resolv-conf.service
[Unit]
Description=Regenerate /etc/resolv.conf to work on WSL in VPN

[Service]
Type=oneshot
ExecStart=%INSTALLDIR%/gen-resolv-conf.sh
